# deployment/med-qa.service
[Unit]
Description=Medical Q&A Portal
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/med-qa-portal/server
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=HOST=localhost
Environment=DOMAIN=localhost

[Install]
WantedBy=multi-user.target